# docker-syncd-vpp/conf/vpp_init.conf
[program:vpp_init.sh]
command=/usr/local/bin/vpp_init.sh
priority=1
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
